# scripts/etc/service/launchCamera/run
#!/bin/bash
export PATH=/usr/local/bin:/usr/bin:/bin
export YARP_CONF=/home/pi/.config/yarp
yarpdev --device opencv_grabber --name /`cat /home/pi/robotName`/img:o
